# deploy/frameshift-backup.service
[Unit]
Description=Stream a FrameShift backup to off-host storage
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=oneshot
User=frameshift
Group=frameshift
ExecStart=/usr/local/sbin/frameshift-backup
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadOnlyPaths=/etc/frameshift /var/lib/frameshift/objects
NoNewPrivileges=true
LockPersonality=true
MemoryDenyWriteExecute=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
SystemCallArchitectures=native
UMask=0077
